# ansible/roles/ready_uwsgi_nginx_web_service/templates/supervisord_uwsgi.conf
{% set uwsgi_workers = ansible_processor_count * 4 %}

[program:{{ app_name }}]
command={{ app_home_dir }}/bin/uwsgi
  --virtualenv {{ venv_dir }}
  --socket {{ app_home_dir }}/run/uwsgi.sock
  --chmod
  --logto {{ app_home_dir }}/log/uwsgi.log
  --log-slow
  --logdate
  --workers {{ uwsgi_workers }}
  --harakiri 30
  --harakiri-verbose
  --max-requests 2000
  --reload-mercy 8
  --cpu-affinity 1
  --master
  --no-orphans
  --vacuum
  --auto-procname
  --procname-prefix "{{ app_name }} "
  --wsgi-file {{ app_src_dir }}/{{ wsgi_file }}
  --callable {{ wsgi_callable }}
user={{ app_user }}
stdout_logfile={{ app_home_dir }}/log/supervisor_uwsgi.log  ; Where to write log messages
redirect_stderr=true                                        ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8             ; Set UTF-8 as default encoding
autostart=true
autorestart=true

; --http 0.0.0.0:5000
; directory={{ app_home_dir }}
; --chdir {{ app_home_dir }}
; --ini-paste production.ini
